# conf.d/nginx-conf.conf
[program:nginx]
command=/app/conf.d/nginx-conf.sh 
autorestart=true