# conf.d/turbovnc.conf
[program:websockify]
command=/opt/TurboVNC/bin/vncserver -noauth 
autorestart=true