# conf.d/vscode.conf
[program:code]
command=code --user-data-dir=""  --no-sandbox  /app/programs
autostart=true