# conf.d/vscode.conf
[program:code]
command=code --user-data-dir="" /app/programs
autostart=true